# configs/launchd/com.local.zfs.automount.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.zfs.automount</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/zfs-automount.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <!-- Intelligent retry: only restart if script exits with error (non-zero) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <!-- Throttle restart attempts to avoid rapid cycling -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/var/log/zfs-automount.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/zfs-automount.log</string>
    <key>UserName</key>
    <string>root</string>
    <key>GroupName</key>
    <string>wheel</string>
    <key>WorkingDirectory</key>
    <string>/var/root</string>
</dict>
</plist>
